# apps/anyscan/anyscan-api.service
[Unit]
Description=AnyScan API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=anyscan
Group=anyscan
WorkingDirectory=/opt/anyscan
EnvironmentFile=-/etc/anyscan/runtime.env
ExecStart=/opt/anyscan/bin/anyscan-api
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/lib/anyscan
StateDirectory=anyscan
RuntimeDirectory=anyscan
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
