# apps/anyscan/anyscan-api.service
[Unit]
Description=AnyScan API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=anyscan
Group=anyscan
WorkingDirectory=/opt/anyscan
EnvironmentFile=-/etc/anyscan/runtime.env
ExecStartPre=+/usr/bin/install -d -m 0755 /opt/anyscan/bin /var/lib/anyscan/build-target
ExecStartPre=+/usr/bin/bash -lc 'TOOLCHAIN=/root/.rustup/toolchains/1.92.0-x86_64-unknown-linux-gnu/bin; export PATH=\"$TOOLCHAIN:/usr/bin:/bin\"; export RUSTC=\"$TOOLCHAIN/rustc\"; cd /root/AnyGPT && \"$TOOLCHAIN/cargo\" build --manifest-path apps/anyscan/Cargo.toml --release --locked --bin anyscan-api --target-dir /var/lib/anyscan/build-target && /usr/bin/install -m 0755 /var/lib/anyscan/build-target/release/anyscan-api /opt/anyscan/bin/anyscan-api'
ExecStart=/opt/anyscan/bin/anyscan-api
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=false
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/lib/anyscan /opt/anyscan/bin
StateDirectory=anyscan
RuntimeDirectory=anyscan
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
