# apps/api/anygpt.service
[Unit]
Description=AnyGPT API Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=anygpt
Group=anygpt
WorkingDirectory=/srv/anygpt/apps/api
TimeoutStartSec=15min
Environment=HOME=/var/lib/anygpt
Environment=NODE_ENV=production
Environment=API_DIST_DIR=dist/production
Environment=PORT=3000
Environment=CLUSTER_WORKERS=6
Environment=PORT_RETRY_COUNT=0
Environment=API_KEYS_FILE=/var/lib/anygpt/keys.json
Environment=API_PROVIDERS_FILE=/var/lib/anygpt/providers.json
Environment=API_MODELS_FILE=/var/lib/anygpt/models.json
Environment=DATA_TMP_DIR=/var/lib/anygpt/tmp
Environment=ANYGPT_LOG_DIR=/var/log/anygpt
Environment=ANYGPT_RUNTIME_DIR=/run/anygpt
ExecStartPre=+/usr/local/bin/bun run build
ExecStart=/usr/local/bin/bun run /srv/anygpt/apps/api/server.launcher.bun.ts
Restart=always
RestartPreventExitStatus=75
RestartSec=5
StateDirectory=anygpt
StateDirectoryMode=0750
LogsDirectory=anygpt
LogsDirectoryMode=0750
RuntimeDirectory=anygpt
RuntimeDirectoryMode=0750
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=yes
ProtectHostname=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectProc=invisible
ProcSubset=pid
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
BindReadOnlyPaths=/root/AnyGPT:/srv/anygpt
BindPaths=/root/AnyGPT/apps/api/dist:/srv/anygpt/apps/api/dist
ReadWritePaths=/srv/anygpt/apps/api/dist /var/lib/anygpt /var/log/anygpt /run/anygpt

[Install]
WantedBy=multi-user.target
Alias=anygpt-api.service
