# apps/api/anygpt-api.service
[Unit]
Description=AnyGPT API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/skullcmd/AnyGPT/apps/api
Environment=NODE_ENV=production
Environment=API_DIST_DIR=dist/production
Environment=PORT=3000
Environment=CLUSTER_WORKERS=1
Environment=PORT_RETRY_COUNT=0
ExecStartPre=/usr/bin/bash -lc '/home/skullcmd/.bun/bin/bun run build'
ExecStart=/usr/bin/bash -lc '/home/skullcmd/.bun/bin/bun run start'
Restart=always
RestartSec=5
User=skullcmd

[Install]
WantedBy=multi-user.target
Alias=anygpt.service
