# apps/api/anygpt-experimental.service
[Unit]
Description=AnyGPT API Service (Experimental)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/skullcmd/AnyGPT/apps/api
Environment=NODE_ENV=experimental
Environment=API_DIST_DIR=dist/experimental
Environment=PORT=3000
Environment=CLUSTER_WORKERS=auto
ExecStartPre=/usr/bin/bash -lc '/home/skullcmd/.bun/bin/bun run build:experimental'
ExecStart=/usr/bin/bash -lc '/home/skullcmd/.bun/bin/bun run start:experimental'
Restart=always
RestartSec=5
User=skullcmd

[Install]
WantedBy=multi-user.target
